# scripts/systemd/inkpi-admin.service
[Unit]
Description=InkPi Local Admin Portal
After=network-online.target inkpi-core.service
Wants=network-online.target inkpi-core.service

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__WORKDIR__
Environment=HOME=__SERVICE_HOME__
Environment=PATH=__SERVICE_PATH__
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-__SERVICE_HOME__/.config/inkpi/admin.env
RuntimeDirectory=inkpi-admin
RuntimeDirectoryMode=0770
ExecStart=__UV_BIN__ run inkpi-admin --host 0.0.0.0 --port 8081 --core-socket /run/inkpi-core/core.sock
Restart=always
RestartSec=5
TimeoutStopSec=30
SuccessExitStatus=143 SIGTERM

[Install]
WantedBy=multi-user.target
